# var/db/repos/ago/app-metrics/queuemetrics/files/queuemetrics.socket
[Unit]
Description=QueueMetrics
Conflicts=queuemetrics.service

[Socket]
ListenStream=8080

[Install]
WantedBy=sockets.target
