# var/db/repos/webos/dev-db/configurator/files/configurator
#!/sbin/runscript

depend() {
	need ls-hubd-pub ls-hubd-prv db8-maindb
}

start() {
	ebegin "Starting ${SVCNAME}"

	start-stop-daemon --make-pidfile --pidfile=/run/configurator.pid --background --start -- configurator -c '{"log":{"appender":{"type":"syslog"},"levels":{"configurator":"notice"}}}' service

	eend 
}

stop() {
	start-stop-daemon --stop -p /run/configurator.pid
}

# vim: ts=4 filetype=gentoo-init-d
