# var/db/repos/webos/dev-db/indb8/files/indb8-maindb
#!/sbin/runscript

depend() {
	need ls-hubd-pub ls-hubd-prv
}

start() {
	ebegin "Starting ${SVCNAME}"

	start-stop-daemon --make-pidfile --pidfile=/run/indb8-maindb.pid --background --start -- /usr/bin/mojodb-luna -c /etc/indb8/maindb.conf

	eend 
}

stop() {
	start-stop-daemon --stop -p /run/indb8-maindb.pid
}

# vim: ts=4 filetype=gentoo-init-d
